# freebsd/atlasmns-trace/pkg-plist
etc/cron.d/atlasmns-tracedataimporter
bin/atlasmns-trace-agent
share/doc/atlasmns-trace/examples/atlasmns-tracedataimporter-configuration
lib/python*/*-packages/AtlasMNS*.egg-info
lib/python*/*-packages/AtlasMNS.py
lib/python*/*-packages/AtlasMNSLogger.py
lib/python*/*-packages/AtlasMNSTools.py
share/doc/atlasmns-trace/examples/NoSQL/README
share/doc/atlasmns-trace/examples/NoSQL/admin.ms
share/doc/atlasmns-trace/examples/NoSQL/database.ms
share/doc/atlasmns-trace/examples/NoSQL/install-database-and-users
share/doc/atlasmns-trace/examples/NoSQL/schema.ms
share/doc/atlasmns-trace/examples/NoSQL/users.ms
share/doc/atlasmns-trace/examples/SQL/README
share/doc/atlasmns-trace/examples/SQL/database.sql
share/doc/atlasmns-trace/examples/SQL/install-database-and-users
share/doc/atlasmns-trace/examples/SQL/schema.sql
share/doc/atlasmns-trace/examples/SQL/users.sql
share/doc/atlasmns-trace/examples/atlasmns-database-configuration
bin/atlasmns-trace-controller
bin/atlasmns-trace-scheduler
man/man1/atlasmns-trace-agent.1.gz
man/man1/atlasmns-trace-controller.1.gz
man/man1/atlasmns-trace-scheduler.1.gz
